# lib/systemd/system/anon-ws-disable-stacked-tor_autogen_port_9151.service
## Copyright (C) 2018 - 2020 ENCRYPTED SUPPORT LP <adrelanos@riseup.net>
## See the file COPYING for copying conditions.

## This file was auto generated by /usr/lib/anon-ws-disable-stacked-tor/systemd-unit-files-generator at package build time.

[Unit]
## IP HARDCODED. See next comment below.
Description=redirect 127.0.0.1:9151 to Whonix-Gateway port 9051
Documentation=https://www.whonix.org/wiki/Dev/anon-ws-disable-stacked-tor
ConditionPathExists=!/run/qubes/this-is-templatevm
## IP HARDCODED. See next comment below.
After=anon-ws-disable-stacked-tor_autogen_port_9151.socket
Requires=anon-ws-disable-stacked-tor_autogen_port_9151.socket

[Service]
## IP HARDCODED. If you want to change IP, set variable GATEWAY_IP through a
## drop-in configuration snippet in /etc/anon-ws-disable-stacked-tor.d
## configuration folder instead.
## Then run the systemd-socket-proxyd systemd-unit-files-generator.
#sudo /usr/lib/anon-ws-disable-stacked-tor/systemd-unit-files-generator
## See Documentation:
## https://www.whonix.org/wiki/Redirect_Whonix-Workstation_Ports_or_Unix_Domain_Socket_Files_to_Whonix-Gateway
ExecStart=/lib/systemd/systemd-socket-proxyd 10.152.152.10:9051
PrivateTmp=yes
