# lib/systemd/system/tor.service.d/50_anon_ws_disable_stacked_tor.conf
## Copyright (C) 2012 - 2023 ENCRYPTED SUPPORT LP <adrelanos@whonix.org>
## See the file COPYING for copying conditions.

[Service]
#### meta start
#### project Whonix
#### category networking and apps
#### workstation_only yes
#### description
## Make 'sudo service tor status' exit '0' for better compatibility.
RemainAfterExit=yes
#### meta end

[Unit]
## Qubes-Whonix:
## Clear 'ConditionPathExists=/run/qubes-service/whonix-gateway' set by
## the qubes-whonix package, which is useful on the gateway but not on the
## workstation.
ConditionPathExists=
