# lib/systemd/system/tor@default.service.d/51_anon_ws_disable_stacked_tor.conf
## Copyright (C) 2017 - 2021 ENCRYPTED SUPPORT LP <adrelanos@whonix.org>
## See the file COPYING for copying conditions.

#### meta start
#### project Whonix
#### category networking and apps
#### workstation_only yes
#### description
## Compatibility with system Tor package.
##
## Overrides systemd unit file by system Tor package for compatibility so
## the dummy Tor binary <code>/usr/bin/tor</code> gets load instead.
##
## Makes all <code>systemctl</code> <code>restart</code> <code>reload</code>
## <code>status</code> commands compatible with dummy Tor.
#### meta end

[Unit]
Description=anon-ws-disable-stacked-tor
Documentation=https://www.whonix.org/wiki/Dev/anon-ws-disable-stacked-tor

[Service]
Type=simple
ExecStartPre=
ExecReload=
AppArmorProfile=
KillSignal=SIGTERM
